# launch-agents/archived/com.personal.ctrld-network-watch.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.personal.ctrld-network-watch</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>/Users/speedybee/bin/ctrld-network-watch.sh</string>
  </array>

  <key>RunAtLoad</key>
  <true/>

  <key>KeepAlive</key>
  <true/>

  <key>StandardOutPath</key>
  <string>/Users/speedybee/Public/Scripts/ctrld-watch.log</string>

  <key>StandardErrorPath</key>
  <string>/Users/speedybee/Public/Scripts/ctrld-watch-error.log</string>
</dict>
</plist>
