# macos/com.abhimehrotra.protondrive-backup.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.abhimehrotra.protondrive-backup</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>/Users/abhimehrotra/Documents/dev/personal-config/scripts/protondrive_backup.sh</string>
    <string>--run</string>
    <string>--no-delete</string>
  </array>

  <key>StartCalendarInterval</key>
  <dict>
    <key>Hour</key>
    <integer>3</integer>
    <key>Minute</key>
    <integer>15</integer>
  </dict>

  <key>StandardOutPath</key>
  <string>/Users/abhimehrotra/Library/Logs/protondrive-backup.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/abhimehrotra/Library/Logs/protondrive-backup.log</string>

  <key>RunAtLoad</key>
  <false/>
</dict>
</plist>
