# media-streaming/archive/launchd/com.abhimehrotra.media.alldebrid.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.abhimehrotra.media.alldebrid</string>

  <!--
    NOTE: launchd does not expand shell variables (e.g. $HOME) in plist
    string values for EnvironmentVariables, StandardOutPath, or
    StandardErrorPath. Only ProgramArguments invoked through `bash -lc`
    will get shell expansion. Replace __HOME__ at install time with a
    concrete path (e.g. via `sed` in your install script) before copying
    this plist to ~/Library/LaunchAgents/.
  -->

  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>-lc</string>
    <string>"$HOME/Library/Media/bin/start-alldebrid.sh"</string>
  </array>

  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    <key>ALD_MOUNT_DIR</key>
    <string>__HOME__/mnt/alldebrid</string>
  </dict>

  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <true/>

  <key>StandardOutPath</key>
  <string>__HOME__/Library/Logs/media/alldebrid.out</string>
  <key>StandardErrorPath</key>
  <string>__HOME__/Library/Logs/media/alldebrid.err</string>
</dict>
</plist>
