# media-streaming/archive/launchd/com.abhimehrotra.media.webdav.plist
<?xml version="1.0" encoding="UTF-8"?>
<!--
  TEMPLATE: launchd does NOT expand shell variables (`$HOME`, `~`) in plist
  fields such as `EnvironmentVariables`, `StandardOutPath`, or
  `StandardErrorPath`. Replace every `__HOME__` placeholder below with the
  absolute path to the target user's home directory (e.g. `/Users/alice`)
  before installing into `~/Library/LaunchAgents`. The `ProgramArguments`
  entry runs through `/bin/bash -lc`, which DOES expand `$HOME` at runtime.
-->
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.abhimehrotra.media.webdav</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>-lc</string>
    <string>$HOME/Library/Media/bin/start-media-server-fast.sh</string>
  </array>

  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    <key>RCLONE_CONFIG</key>
    <string>__HOME__/.config/rclone/rclone.conf</string>
    <key>MEDIA_CACHE_DIR</key>
    <string>__HOME__/Library/Application Support/MediaCache</string>
    <key>MEDIA_CREDENTIALS_FILE</key>
    <string>__HOME__/.config/media-server/credentials</string>
  </dict>

  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <true/>

  <key>StandardOutPath</key>
  <string>__HOME__/Library/Logs/media/webdav.out</string>
  <key>StandardErrorPath</key>
  <string>__HOME__/Library/Logs/media/webdav.err</string>
</dict>
</plist>
