# maintenance/launchd/com.abhimehrotra.maintenance.servicemonitor.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.abhimehrotra.maintenance.servicemonitor</string>
    
    <key>ProgramArguments</key>
    <array>
        <string>/Users/abhimehrotra/Documents/dev/personal-config/maintenance/bin/service_monitor.sh</string>
    </array>
    
    <key>EnvironmentVariables</key>
    <dict>
        <key>AUTOMATED_RUN</key>
        <string>1</string>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    </dict>
    
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>8</integer>
        <key>Minute</key>
        <integer>35</integer>
    </dict>
    
    <key>StandardOutPath</key>
    <string>/Users/abhimehrotra/Library/Logs/maintenance/servicemonitor-stdout.log</string>
    
    <key>StandardErrorPath</key>
    <string>/Users/abhimehrotra/Library/Logs/maintenance/servicemonitor-stderr.log</string>
    
    <key>RunAtLoad</key>
    <false/>
    
    <key>Nice</key>
    <integer>5</integer>
    
    <key>ProcessType</key>
    <string>Background</string>
</dict>
</plist>
