# maintenance/launchd/com.abhimehrotra.maintenance.screencapture-nag-remover.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.abhimehrotra.maintenance.screencapture-nag-remover</string>
    
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/opt/homebrew/sbin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
        <key>MAINTENANCE_HOME</key>
        <string>/Users/speedybee/Library/Maintenance</string>
    </dict>
    
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/speedybee/Library/Maintenance/bin/screencapture_nag_remover.sh</string>
    </array>
    
    <key>StandardOutPath</key>
    <string>/Users/speedybee/Library/Logs/maintenance/screencapture_nag_remover.out</string>
    
    <key>StandardErrorPath</key>
    <string>/Users/speedybee/Library/Logs/maintenance/screencapture_nag_remover.err</string>
    
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>10</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>
    
    <key>RunAtLoad</key>
    <true/>
</dict>
</plist>
